# etc/crontab
# Updates the Let's Encrypt Certificates
@daily 						/usr/local/bin/letsencrypt renew
